# telegram-jellyfin-bot/services/telegram-media-bot.service
[Unit]
Description=Telegram Userbot Media Bot
After=network.target

[Service]
Type=simple
User=YOUR_USER
WorkingDirectory=/opt/telegram-jellyfin-bot
ExecStart=/opt/telegram-jellyfin-bot/venv/bin/python bot.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target